# Docker/sync/cron.d/sync
MAILTO=""
* * * * * /bin/date >> /tmp/hello.txt

