# Docker/sync/cron.d/sync
* * * * * /bin/sh /usr/src/sync/scripts/sync.sh

